# systemd/hosaka-field-terminal-headless.service
[Unit]
Description=Hosaka Field Terminal Appliance (Headless)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hosaka-field-terminal
Environment=HOSAKA_STATE_PATH=/var/lib/hosaka/state.json
Environment=PYTHONUNBUFFERED=1
Environment=HOSAKA_BOOT_MODE=headless
ExecStart=/opt/hosaka-field-terminal/.venv/bin/python -m hosaka
Restart=always
RestartSec=2
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
