# systemd/hosaka-field-terminal.service
[Unit]
Description=Hosaka Field Terminal Appliance
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hosaka-field-terminal
Environment=HOSAKA_STATE_PATH=/var/lib/hosaka/state.json
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/hosaka-field-terminal/.venv/bin/python -m hosaka
Restart=always
RestartSec=2
StandardInput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

[Install]
WantedBy=multi-user.target
